# raspi/userland/target/brummbeere.target
[Unit]
Description=Brummbeere application
Wants=multi-user.target network-online.target
After=multi-user.target network-online.target
