# raspi/userland/target/brummbeere.service
[Unit]
Description=Brummbeere
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/brummbeere.sh

[Install]
WantedBy=brummbeere.target
